# scripts/systemd/autopilot-metrics.service
[Unit]
Description=Autopilot Weekly Metrics Report Generator
Documentation=https://github.com/openagents/openagents
After=network.target

[Service]
Type=oneshot
User=christopherdavid
WorkingDirectory=/home/christopherdavid/code/openagents
ExecStart=/home/christopherdavid/code/openagents/scripts/weekly-metrics-report.sh

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/christopherdavid/.cargo/bin"
Environment="AUTOPILOT_DB=/home/christopherdavid/code/openagents/autopilot.db"
Environment="METRICS_OUTPUT=/home/christopherdavid/code/openagents/docs/metrics/weekly"

# Optional: Send notifications via webhook
# Environment="NOTIFY_WEBHOOK=https://hooks.slack.com/services/YOUR/WEBHOOK/URL"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autopilot-metrics

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
