# scripts/systemd/autopilot-metrics.timer
[Unit]
Description=Autopilot Weekly Metrics Report Timer
Documentation=https://github.com/openagents/openagents
Requires=autopilot-metrics.service

[Timer]
# Run every Monday at 9:00 AM
OnCalendar=Mon *-*-* 09:00:00

# If system was off when timer should have run, run it on next boot
Persistent=true

# Randomize start time by up to 5 minutes to avoid thundering herd
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
