# docs/autopilot/weekly-report.timer
[Unit]
Description=OpenAgents Autopilot Weekly Report Timer
Requires=weekly-report.service

[Timer]
# Run every Monday at 9:00 AM
OnCalendar=Mon *-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
