# systemd/autopilot-baseline-update.service
[Unit]
Description=Autopilot Weekly Baseline Update
Documentation=https://github.com/OpenAgentsInc/openagents
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/code/openagents
ExecStart=%h/.cargo/bin/openagents autopilot baseline update
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
