# systemd/autopilot-baseline-update.timer
[Unit]
Description=Autopilot Weekly Baseline Update Timer
Documentation=https://github.com/OpenAgentsInc/openagents
Requires=autopilot-baseline-update.service

[Timer]
# Run every Monday at 00:00
OnCalendar=Mon *-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
