# docs/autopilotd.service
[Unit]
Description=Autopilot Daemon - AI worker supervisor
Documentation=file:///home/christopherdavid/code/openagents/AGENTS.md
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/rm -f /home/christopherdavid/.autopilot/autopilotd.pid /home/christopherdavid/.autopilot/autopilotd.sock
ExecStart=/home/christopherdavid/.cargo/bin/cargo run --manifest-path /home/christopherdavid/code/openagents/Cargo.toml -p autopilot --bin autopilotd -- --workdir /home/christopherdavid/code/openagents --project openagents
WorkingDirectory=/home/christopherdavid/code/openagents
Restart=always
RestartSec=5
# Clean up stale files on stop
ExecStopPost=/bin/rm -f /home/christopherdavid/.autopilot/autopilotd.pid /home/christopherdavid/.autopilot/autopilotd.sock

# Environment
Environment=PATH=/home/christopherdavid/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=default.target
